# supervisord.conf
[supervisord]
nodaemon=true

[program:php7.3-fpm]
command=/usr/sbin/php-fpm7.3 --nodaemonize --fpm-config=/etc/php/7.3/fpm/pool.d/www.conf
numprocs=1
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
numprocs=1
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:freeswitch]
command=/usr/bin/freeswitch -u www-data -g www-data -ncwait -nonat
numprocs=1
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
